# files/systemd/kiosk.service
[Unit]
Description=Kiosk Autostart

[Install]
WantedBy=xdg-desktop-autostart.target
